# systemd/firefox-updater.service
[Unit]
Description=Update Firefox Developer edition

[Service]
Type=oneshot
ExecStart={{SCRIPT_LOC}}
